# supervisord.conf
[supervisord]
nodaemon=true

[program:python_server]
command=python main.py
directory=/app
autorestart=true